# crontab
# Job IDs must match trmx casting catalog (server/utils/castingCronCatalog.ts).
# TRMX_CRON_HEARTBEAT_* / CRON_HEARTBEAT_SECRET: set in Coolify (see docker-entrypoint → /etc/cron-env.sh).

# Project casting — 10:00 and 18:00 PT (CRON_TZ in container)
0 18 * * * cd /app && ./scripts/cron-with-heartbeat.sh project_casting_18 npm run start:project_casting
0 10 * * * cd /app && ./scripts/cron-with-heartbeat.sh project_casting_10 npm run start:project_casting

# Auditions free — 18:00
0 18 * * * cd /app && ./scripts/cron-with-heartbeat.sh auditions_free npm run start:auditions_free

# All castings — 18:00
0 18 * * * cd /app && ./scripts/cron-with-heartbeat.sh all_castings npm run start:all_castings

# Meilisearch — every 10 minutes
*/10 * * * * cd /app && ./scripts/cron-with-heartbeat.sh melisearch_projects npm run start:add_projects_to_melisearch

# Project markets — every 10 minutes
*/10 * * * * cd /app && ./scripts/cron-with-heartbeat.sh project_markets npm run update:project_markets
